# templates/default/descartes-net.conf.erb
description "descartes network"

start on starting networking
stop on stopped networking

env USE_LXC_BRIDGE="true"
env LXC_BRIDGE="lxcbr0"
env LXC_ADDR="<%= @addr %>"
env LXC_NETMASK="255.255.255.0"
env LXC_NETWORK="<%= @network %>"
env varrun="/var/run/descartes-net"

pre-start script
     [ -f "${varrun}/network_up" ] && exit 0;
    if [ -d /sys/class/net/${LXC_BRIDGE} ]; then
	if [ ! -f ${varrun}/network_up ]; then
	    # bridge exists, but we didn't start it
	    stop;
	fi
	exit 0;
    fi

	# set up the lxc network
    echo 1 > /proc/sys/net/ipv4/ip_forward
    mkdir -p ${varrun}
    brctl addbr ${LXC_BRIDGE}
    ifconfig ${LXC_BRIDGE} ${LXC_ADDR} netmask ${LXC_NETMASK} up
#    iptables -t nat -A POSTROUTING -s ${LXC_NETWORK} ! -d ${LXC_NETWORK} -j MASQUERADE
    touch ${varrun}/network_up
end script

post-stop script
    [ -f "${varrun}/network_up" ] || exit 0;
    # if $LXC_BRIDGE has attached interfaces, don't shut it down
    ls /sys/class/net/${LXC_BRIDGE}/brif/* > /dev/null 2>&1 && exit 0;

    if [ -d /sys/class/net/${LXC_BRIDGE} ]; then
	ifconfig ${LXC_BRIDGE} down
	#iptables -t nat -D POSTROUTING -s ${LXC_NETWORK} -j MASQUERADE || true
	brctl delbr ${LXC_BRIDGE}
    fi
    rm -f ${varrun}/network_up
end script
